# texture.service
[Unit]
Description=Texture ML Learning Studio
After=network.target

[Service]
Type=forking
User=YOUR_USERNAME
WorkingDirectory=/path/to/texture
ExecStart=/usr/bin/make prod
ExecStop=/usr/bin/make prod-stop
Restart=on-failure
RestartSec=10

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
